# install/etc/profile.d/py_venv.sh
source /opt/venv/bin/activate
source /opt/bin/setenv.sh
